# docker-compose/supervisord.conf
[supervisord]
nodaemon=true

[program:app-nginx]
command=nginx
