# docker-compose/supervisord.conf
[supervisord]
nodaemon=true

[program:app-uwsgi]
command=uwsgi --ini ./uwsgi.ini

[program:app-nginx]
command=nginx
